# deploy/systemd-user/tmuxgo-gateway.service
[Unit]
Description=TmuxGo Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__TMUXGO_ROOT__
ExecStart=__TMUXGO_NPM__ run --workspace=gateway start
Restart=always
RestartSec=3
KillMode=control-group
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=PATH=__TMUXGO_PATH__

[Install]
WantedBy=tmuxgo.target
