# deploy/systemd-user/tmuxgo.target
[Unit]
Description=TmuxGo Suite
Wants=tmuxgo-gateway.service tmuxgo-frontend.service tmuxgo-agent.service

[Install]
WantedBy=default.target
