# deploy/systemd-user/tmuxgo-agent.service
[Unit]
Description=TmuxGo Agent
After=network-online.target tmuxgo-gateway.service
Wants=network-online.target tmuxgo-gateway.service

[Service]
Type=simple
WorkingDirectory=__TMUXGO_ROOT__
ExecStart=__TMUXGO_NPM__ run --workspace=agent start
Restart=always
RestartSec=3
KillMode=process
Environment=NODE_ENV=production
Environment=GATEWAY_URL=ws://127.0.0.1:3001/api/stream
Environment=HOST_ID=agent-local
Environment=HOST_NAME=%H
Environment=PATH=__TMUXGO_PATH__

[Install]
WantedBy=tmuxgo.target
